# files/openwrt/etc/init.d/tg-ws-proxy
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=66

NAME="tg-ws-proxy"
PROG="/usr/bin/tg-ws-proxy"
CONFIG_DIR="/etc/tg-ws-proxy"
CONFIG_FILE="$CONFIG_DIR/config.conf"
SECRET_FILE="$CONFIG_DIR/secret.conf"
LOGFILE="/var/log/tg-ws-proxy.log"

load_config() {
        [ -f "$CONFIG_FILE" ] || {
                echo "Config file not found: $CONFIG_FILE"
                return 1
        }

        [ -f "$SECRET_FILE" ] || {
                echo "Secret file not found: $SECRET_FILE"
                return 1
        }

        . "$CONFIG_FILE"
        . "$SECRET_FILE"

        [ -n "$SECRET" ] || {
                echo "SECRET is empty in $SECRET_FILE"
                return 1
        }

        [ -n "${DC_IP_DEFAULT_POOL+x}" ] || DC_IP_DEFAULT_POOL=""
        [ -n "${CFPROXY_FALLBACK_DOMAINS+x}" ] || CFPROXY_FALLBACK_DOMAINS=""
        [ -n "${CFPROXY_DOMAINS_URL+x}" ] || CFPROXY_DOMAINS_URL=""
        [ -n "${FAKE_TLS_DOMAIN+x}" ] || FAKE_TLS_DOMAIN=""
        return 0
}

print_link() {
        local link_host br_ip

        load_config || return 1
        [ -n "$HOST" ] && [ -n "$PORT" ] && [ -n "$SECRET" ] || return 1

        link_host="$HOST"
        if [ "$link_host" = "0.0.0.0" ]; then
                br_ip="$(ip -f inet addr show dev br-lan 2>/dev/null | sed -n 's/.*inet \([0-9.]\+\)\/.*/\1/p' | head -n 1)"
                [ -n "$br_ip" ] || br_ip="$(ip -f inet addr show dev br0 2>/dev/null | sed -n 's/.*inet \([0-9.]\+\)\/.*/\1/p' | head -n 1)"
                [ -n "$br_ip" ] && link_host="$br_ip"
        fi

        if [ -n "$FAKE_TLS_DOMAIN" ]; then
                domain_hex="$(printf '%s' "$FAKE_TLS_DOMAIN" | od -An -tx1 2>/dev/null | tr -d ' \n')"
                echo "Connect link: tg://proxy?server=$link_host&port=$PORT&secret=ee$SECRET$domain_hex"
        else
                echo "Connect link: tg://proxy?server=$link_host&port=$PORT&secret=dd$SECRET"
        fi
}

start_service() {
        config_load "$NAME"

        local enabled user
        config_get_bool enabled "main" "enabled" "0"
        [ "$enabled" -eq "1" ] || return 0

        load_config || return 1

        config_get user "main" "user" "root"

        procd_open_instance "$NAME.main"
        procd_set_param command "$PROG" \
                --host "$HOST" \
                --port "$PORT" \
                --secret "$SECRET" \
                --dc-ip-default "$DC_IP_DEFAULT" \
                --dc-ip-default-pool "$DC_IP_DEFAULT_POOL"

        if [ -n "$CFPROXY_FALLBACK_DOMAINS" ]; then
                procd_append_param command --cfproxy-domains "$CFPROXY_FALLBACK_DOMAINS"
        fi
        if [ -n "$CFPROXY_DOMAINS_URL" ]; then
                procd_append_param command --cfproxy-domains-url "$CFPROXY_DOMAINS_URL"
        fi
        if [ -n "$FAKE_TLS_DOMAIN" ]; then
                procd_append_param command --fake-tls-domain "$FAKE_TLS_DOMAIN"
        fi

        if [ -n "$EXTRA_ARGS" ]; then
                set -- $EXTRA_ARGS
                while [ "$#" -gt 0 ]; do
                        procd_append_param command "$1"
                        shift
                done
        fi

        if [ "$LOG_LEVEL" = "1" ]; then
                procd_append_param command --log-file "$LOGFILE"
        fi

        procd_set_param stdout 1
        procd_set_param stderr 1
        procd_set_param user "$user"
        procd_set_param respawn
        procd_close_instance
        print_link
        return 0
}

status_service() {
        if [ -n "$(pidof tg-ws-proxy 2>/dev/null)" ]; then
                echo "tg-ws-proxy is running"
                print_link
                return 0
        fi

        echo "tg-ws-proxy is not running"
        return 1
}
